# docs/brewmbled.service.sample
[Unit]
Description=Brewmble Daemon
After=network-online.target
Wants=network-online.target
Documentation=https://github.com/hebra/brewmble

[Service]
Type=simple
ExecStart=/opt/brewmble/.cargo/bin/brewmbled
Restart=on-failure
User=brewmble

# Configuration via environment variables.
# Uncomment and set these as needed for your environment.

# Port to listen on (default: 8080).
# Environment=BREWMBLE_DAEMON_PORT=8080

# Hostname for mDNS registration (default: system hostname).
# Environment=BREWMBLE_DAEMON_HOSTNAME=my-node

# Explicit IP address for mDNS registration.
# Environment=BREWMBLE_DAEMON_IP=192.168.1.10

# API key for authentication. If not set, a random one is generated on startup.
# Environment=BREWMBLE_DAEMON_API_KEY=your-secure-api-key

# Logging level (error, warn, info, debug, trace).
# Environment=RUST_LOG=info

# Interval in minutes between 'apt-get update' calls (Linux only).
# Default is 360 (6 hours). Set to 0 to force an update on every status check.
# Environment=BREWMBLE_APT_UPDATE_INTERVAL=360

[Install]
WantedBy=multi-user.target
